# systemd/bor-resync.timer
[Unit]
Description=Hourly timer for browser-on-ram
BindsTo=bor.service

[Timer]
OnUnitActiveSec=1h

# vim: ft=systemd
