# systemd/user/bor-resync.timer
[Unit]
Description=Hourly timer for browser-on-ram
BindsTo=bor.service

[Timer]
OnActiveSec=1h
OnUnitActiveSec=1h

# vim: ft=systemd
