# systemd/user/bor-resync.service
[Unit]
Description=Timed resync
BindsTo=bor.service
After=bor.service

[Service]
Type=oneshot
ExecStart=bor --resync --verbose

# vim: ft=systemd
