# systemd/user/bor-resync.service
[Unit]
Description=Timed resync

[Service]
Type=oneshot
ExecStart=bor --resync --verbose
Slice=background.slice

# vim: ft=systemd
