# systemd/user/bor-sleep.target
[Unit]
Description=System sleep target for browser-on-ram
StopWhenUnneeded=yes
